# src/systemd/birdnet-power.service
[Unit]
Description=BirdNET-Pi Remix Power Manager
After=multi-user.target

[Service]
User=root
Environment=PYTHONPATH=/opt/birdnet-remix/src
ExecStart=/opt/birdnet-remix/venv/bin/python3 /opt/birdnet-remix/src/power/daemon.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=birdnet-remix.target
